# scripts/systemd/bwc-restore-drill.service
[Unit]
Description=BWC weekly backup-restore drill (decrypt + restore latest .sql.age to scratch DB)
# Without OnFailure, a drill detecting backup corruption or restore
# breakage would log to journal but never alert the operator — defeating
# the drill's purpose.
OnFailure=bwc-alert@%n.service

[Service]
Type=oneshot
# Runs as root: needs CREATE DATABASE / DROP DATABASE privilege that
# bwc_migrator does not have (setup.sh's bwc_migrator GRANT block
# scopes it to bwc.* only), AND needs read access to
# /etc/bwc/age-identity (mode 0600 root:root). Root-via-socket auth
# to MariaDB avoids any in-script credentials.
User=root
# Drill workload is mysql restore + age decrypt streaming — share the
# same I/O priority class as the daily backups. Most relevant when
# the worst-case drill end (Sun 02:00 + 30min jitter + 30min cap
# ~ 03:00) just-touches bwc-db-backup's 03:00 start window; idle class
# yields to db-backup's own IO if they overlap. Caveat: this only
# moderates the drill's own bash/age/gunzip IO; mysqld's writes to
# /var/lib/mysql happen in a separate process tree and inherit
# mysqld's own (default) IO class.
IOSchedulingClass=idle
# Leading `-` makes systemd skip the chdir if the symlink is missing
# (pre-first-deploy state) instead of failing with status=200/CHDIR
# and triggering OnFailure spam on every timer tick.
WorkingDirectory=-/opt/bwc/current
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
# Reads /backup/bwc/db and /etc/bwc/age-identity, writes only the
# success marker /var/lib/bwc/last-restore-drill.ok and the lock file
# /var/lib/bwc/.bwc-restore-drill.lock. The scratch DB lives in the
# MariaDB data directory (/var/lib/mysql, owned by mysql:mysql) —
# accessed via the unix socket, NOT via filesystem write, so no
# ReadWritePaths entry is needed for it.
ReadWritePaths=/var/log/bwc /var/lib/bwc
# mysql client + age decrypt + gunzip streaming pipeline; 512M caps
# the worst-case dump size we'd reasonably restore (BWC is small —
# bound by mysqldump buffer + gunzip dictionary, not by total dump size
# since we stream, not buffer in memory). MemorySwapMax=0 forbids
# swap-spillover so a runaway restore can't push the host into thrash.
MemoryMax=512M
MemorySwapMax=0
TasksMax=128
# Cap total runtime. Drill restores typically <5min; 30m ceiling
# absorbs slow-disk + cleanup_scratch on a transient mysql stall.
# Without an override, systemd's default 90s would kill the drill
# mid-restore on any nontrivial dump.
TimeoutStartSec=1800s
# Short-circuit when the release artifact hasn't landed yet — the
# ExecStart path lives in the symlinked current release. Without this,
# the timer would spam OnFailure on pre-first-deploy ticks.
ExecCondition=/usr/bin/test -x /opt/bwc/current/scripts/restore-drill.sh
ExecStart=/opt/bwc/current/scripts/restore-drill.sh
