# scripts/systemd/bwc-static-prune.timer
[Unit]
Description=BWC static-pool pruner timer

[Timer]
# 05:30 UTC — between media-verify and the first traffic of the day.
# Pruning removes /_next/static/* chunks no longer referenced by
# any release on disk, capping the static pool's footprint.
OnCalendar=*-*-* 05:30:00
RandomizedDelaySec=10min
Persistent=true
Unit=bwc-static-prune.service

[Install]
WantedBy=timers.target
