# scripts/systemd/cavecms-alert@.service
[Unit]
Description=CaveCMS OnFailure alert for %i

[Service]
Type=oneshot
# Dedicated alert env file — contains ONLY ADMIN_ALERT_EMAIL, NOT
# the full production env. Prevents mailx (and any subprocess it
# spawns) from inheriting JWT_SECRET / DATABASE_URL / SMTP_PASS.
EnvironmentFile=/etc/cavecms/alert.env
# Cap the unit at 60s so a hung MTA can't pin the alert oneshot
# forever (default TimeoutStartSec is 90s; tightening to 60s
# matches mailx's default smtp-connect-timeout).
TimeoutStartSec=60s
# Sandboxing — this template is the ONE unit that runs after any
# other service fails, so it has the broadest reach in the bundle
# and deserves the strongest containment. CAP_DAC_READ_SEARCH lets
# journalctl read /var/log/journal regardless of FS perms;
# everything else is locked down. ReadWritePaths covers the local
# MTA spool (postfix/sendmail typically write under /var/spool).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=
ReadWritePaths=/var/spool /var/log
# journalctl + mailx are streaming; 128M is generous.
MemoryMax=128M
MemorySwapMax=0
TasksMax=64
# Three protections vs. silent alert loss:
# 1. `logger -t cavecms-alert "FAIL $1"` always writes to syslog, so
#    `journalctl -t cavecms-alert` shows the trail even if mail is broken.
# 2. `mailx -- "$ADMIN_ALERT_EMAIL"` terminates option parsing — a
#    recipient that starts with `-` cannot inject `-S` flags.
# 3. If mailx fails (no MTA, queue full, DNS down), a second logger
#    line records the delivery failure. Operator scans
#    `journalctl -t cavecms-alert | grep "MAIL FAILED"` to detect.
# %I (capital) is the fully-unescaped instance name; the bash
# positional-arg pattern (`_ %I`) prevents shell injection from
# unit names containing systemd-escape sequences.
ExecStart=/bin/bash -c 'msg=$(journalctl -u "$1" --no-pager -n 50); logger -t cavecms-alert "FAIL $1"; printf "%s\n" "$msg" | mailx -s "[cavecms] $1 failed" -- "${ADMIN_ALERT_EMAIL:-root}" || logger -t cavecms-alert "MAIL FAILED for $1"' _ %I
