# scripts/systemd/cavecms-db-backup.service
[Unit]
Description=CaveCMS nightly database backup
OnFailure=cavecms-alert@%n.service

[Service]
Type=oneshot
User=cavecms
Group=backup
# Group=backup overrides the cavecms user's default primary group AND drops
# all supplementary groups. The cavecmsstate group scopes /var/lib/cavecms rwx
# (lock + marker files) to cavecms-user services only — NOT to www-data
# which is in group `cavecms` for nginx serving uploads. See setup.sh
# cavecmsstate-group block for the full rationale.
SupplementaryGroups=cavecmsstate
# Backups are bulk I/O; deprioritise so the runtime process isn't
# starved.
IOSchedulingClass=idle
EnvironmentFile=/etc/cavecms/env.production
# Leading `-` makes systemd skip the chdir if the symlink is
# missing (pre-first-deploy state) instead of failing
# `status=200/CHDIR` and triggering OnFailure spam on every timer
# tick. ExecCondition below then cleanly short-circuits with rc=0.
WorkingDirectory=-/opt/cavecms/current
# Defence-in-depth sandboxing. Even if db-backup.sh is compromised,
# the process can't read /home, /root, or arbitrary FS; can't gain
# new privileges via setuid binaries; capability set is empty so no
# `mount`, no `kill -9 root processes`, no raw socket. Writes are
# restricted to the dirs it needs.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
ReadWritePaths=/backup /var/log/cavecms /var/lib/cavecms
# Resource bounds — mysqldump streams from a transactional snapshot
# so RSS stays low; 256M is generous. Without bounds, a runaway
# dump could OOM-kill the PM2 worker (which usually has the
# highest oom_score on a 1 GB droplet).
MemoryMax=256M
MemorySwapMax=0
TasksMax=256
# Cap total runtime. mysqldump on a healthy small DB takes <30s; the
# 30m ceiling absorbs slow-disk + retry budget. Without an override,
# systemd's default 90s kills legitimate large-DB dumps mid-pipeline.
TimeoutStartSec=1800s
ExecCondition=/usr/bin/test -x /opt/cavecms/current/scripts/db-backup.sh
ExecStart=/opt/cavecms/current/scripts/db-backup.sh
