# scripts/systemd/cavecms-restore-drill.timer
[Unit]
Description=CaveCMS restore drill timer

[Timer]
# Weekly cadence, Sunday at 02:00 local. Placed BEFORE the other
# Sunday timers (db-backup 03:00, uploads-backup 04:30, media-verify
# 05:00, static-prune 05:30, cron-purge 06:00) so the drill's mysql
# restore IO doesn't contend with the day's first backup. Worst-case
# drill end (02:00 + 30m jitter + 30m TimeoutStartSec ~ 03:00) can
# just-touch db-backup's start window, but they hit different
# resources (drill writes /var/lib/mysql, db-backup reads cavecms DB +
# writes /backup) so contention is minimal. An earlier iteration of
# this file scheduled the drill at Sun 04:00, which collided with
# uploads-backup at 04:30 (both writing /backup) — moved to 02:00 in
# the cluster-3 cohesive audit.
#
# RandomizedDelaySec=30min scatters across a fleet (single-droplet
# right now, multi-host-future safe). Persistent=true catches up a
# missed run after a downtime window — the drill should still execute
# when the host comes back up even if the scheduled time has passed.
OnCalendar=Sun *-*-* 02:00:00
RandomizedDelaySec=30min
Persistent=true
Unit=cavecms-restore-drill.service

[Install]
WantedBy=timers.target
