# scripts/systemd/cavecms-cron-crm-retry.timer
[Unit]
Description=CaveCMS CRM dispatch retry timer

[Timer]
# Every minute. The retry backoff buckets are 5min / 30min / 4h /
# 24h — a 1-minute polling interval keeps retry latency under
# (backoff bucket + 1min). Persistent=true so a missed tick (host
# down, service reload) doesn't drop a retry-due row indefinitely.
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=15s
Persistent=true
Unit=cavecms-cron-crm-retry.service

[Install]
WantedBy=timers.target
