# scripts/systemd/cavecms-cron-purge.service
[Unit]
Description=CaveCMS nightly retention purge
# Note: there is NO systemd-level dependency on the backup units.
# An `After=` clause would only order units inside a SINGLE
# transaction — since the timers fire independently from the
# system clock, `After=` provides no real safety. Wall-clock
# spacing (uploads-backup 04:30 → cron-purge 06:00) is the
# scheduling buffer; cron-purge.ts itself reads a marker file
# (`/var/lib/cavecms/last-uploads-backup.ok` — the HEALTH marker, touched
# by uploads-backup.sh on EVERY successful exit INCLUDING no-changes
# skips) written by the backup script on success and aborts if the
# marker is older than 24h. NOTE: do NOT confuse this with the hidden
# `/var/lib/cavecms/.last-uploads-backup.scan` SCAN marker, whose mtime
# tracks tar-start time and is consumed only by uploads-backup.sh's
# own find-newer gate — reading it would produce a misleading age
# signal on no-changes days.
OnFailure=cavecms-alert@%n.service

[Service]
Type=oneshot
# 30-minute hard ceiling. The PR-4 cron-purge HTTP loopback to
# /api/internal/revalidate-tags has no per-request timeout in undici by
# default, so a wedged Next pool would otherwise hang this oneshot in
# `activating` indefinitely — silencing the systemd OnFailure= alert AND
# suppressing the next daily timer fire (oneshot is serialized). 30 min
# is far above the realistic upper bound for any worst-case backlog
# sweep + media unlink; past that, systemd hard-kills and OnFailure= fires.
TimeoutStartSec=30min
User=cavecms
Group=cavecms
# Group=cavecms (explicit, same as the user's default primary) drops all
# supplementary groups; add cavecmsstate so cron-purge can read the marker
# files in /var/lib/cavecms (root:cavecmsstate 2770, files inherit group cavecmsstate
# via setgid). See setup.sh cavecmsstate-group block for the rationale.
SupplementaryGroups=cavecmsstate
# Belt-and-suspenders: cron-purge.ts hard-refuses to run unless
# NODE_ENV=production. Setting it directly on the unit (BEFORE the
# EnvironmentFile load) guarantees the gate passes even if the file
# is missing the line or out-of-date — EnvironmentFile entries override
# Environment= per systemd.exec(5), so if env.production sets NODE_ENV
# it still wins; this is purely a floor, not a ceiling.
Environment=NODE_ENV=production
EnvironmentFile=/etc/cavecms/env.production
# `-` prefix tolerates missing symlink pre-first-deploy.
WorkingDirectory=-/opt/cavecms/current
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
# Hard-deletes unreferenced upload files, so /opt/cavecms/uploads is
# read-write here even though most timer services have it read-only.
ReadWritePaths=/opt/cavecms/uploads /var/log/cavecms /var/lib/cavecms
# tsx + Node + Drizzle baseline is ~200M; the purge sweeps batches
# of soft-deleted rows so memory should stay flat. 768M ceiling
# leaves headroom for unexpectedly large delete sets without
# starving the PM2 worker.
MemoryMax=768M
MemorySwapMax=0
TasksMax=256
# Absolute path to tsx avoids relying on systemd's default PATH or
# on `tsx` being installed globally. tsx is shipped as a runtime
# dependency in the release artifact (package.json moved from
# devDependencies in cluster 1 cleanup).
ExecCondition=/usr/bin/test -x /opt/cavecms/current/node_modules/.bin/tsx
ExecStart=/opt/cavecms/current/node_modules/.bin/tsx scripts/cron-purge.ts
